# ddp-receiver.service.in
# Generated by CMake. Reconfigure the project instead of editing this file.
# Panel geometry: @DDP_MATRIX_PANEL_COLS@x@DDP_MATRIX_PANEL_ROWS@ pixels per panel
# Panels: total=@DDP_MATRIX_PANEL_COUNT@ (chain=@DDP_MATRIX_PANEL_CHAIN@, parallel=@DDP_MATRIX_PANEL_PARALLEL@)
# Canvas: @DDP_MATRIX_CANVAS_WIDTH@x@DDP_MATRIX_CANVAS_HEIGHT@
# ExecStart args: @DDP_MATRIX_SERVICE_ARGS_STRING@

[Unit]
Description=DDP Receiver LED Matrix Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=@DDP_MATRIX_INSTALL_DIR@/ddp-receiver @DDP_MATRIX_SERVICE_ARGS_STRING@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
